# templates/cron.d/nginx-cloudwatch.erb
*/5 * * * * root /etc/<%= @brand_name %>/nginx-aws-logger.py
*/5 * * * * root /etc/<%= @brand_name %>/mon-put-instance-data.pl --mem-util --mem-avail --disk-space-util --disk-space-avail --disk-path=/ --from-cron
